# imm/buildroot/files/etc/profile.d/shell.sh
export EDITOR="nano"
export PAGER="less"
alias bridge="bridge -color=auto"
alias diff="diff --color=auto"
alias grep="grep --color=auto"
alias ip="ip -color=auto"