# ax6/buildroot/files/etc/uci-defaults/99-qol_fixes
                uci set wireless.radio0.country='US'
                uci set wireless.radio1.country='US'
                uci set wireless.radio2.country='US'
                uci set wireless.radio1.disabled=0
                uci set wireless.radio2.disabled=0
                uci set pbuf.opt.memory_profile=auto
                uci set network.globals.packet_steering=0
                uci set firewall.@defaults[0].flow_offloading=0
                uci set ecm.@general[0].enable_bridge_filtering=0
                uci set system.@system[0].hostname="Router"
                uci set system.@system[0].timezone="CST-8"
                uci set system.@system[0].zonename="Asia/Shanghai"
                uci delete system.ntp.server
                uci add_list system.ntp.server="ntp.aliyun.com"
                uci add_list system.ntp.server="ntp.tencent.com"
                uci add_list system.ntp.server="ntp.ntsc.ac.cn"
                uci add_list system.ntp.server="cn.ntp.org.cn"
                uci set luci.main.lang='zh_cn'
                uci set luci.main.mediaurlbase='/luci-static/bootstrap'
                uci commit